# puppetmaster/modules/config/templates/apps/unisubs/upstart.celeryd.conf.erb
# unisubs celeryd

description "unisubs: celeryd (<%= @env %>)"
start on runlevel [2345]
stop on runlevel [06]

exec <%= @ve_dir %>/bin/python <%= @project_dir %>/manage.py <% if @env == 'local' %>--purge <% end %>celeryd -B -c 6 -E -S djcelery.schedulers.DatabaseScheduler -f /var/log/celeryd.<%= @env %>.log --settings <%= @settings_module %>
