# puppetmaster/modules/config/templates/apps/unisubs/upstart.celerycam.conf.erb
# unisubs celerycam

description "unisubs: celerycam (<%= @env %>)"
start on runlevel [2345]
stop on runlevel [06]

exec <%= @ve_dir %>/bin/python <%= @project_dir %>/manage.py celerycam -f /var/log/celerycam.<%= @env %>.log --settings <%= @settings_module %>
